# cloud/supervisor.d/mosquitto.ini
[program:mosquitto]
priority=20
directory=/tmp
command=/usr/sbin/mosquitto -c /mosquitto/config/mosquitto.conf
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

